# supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
command=python bot/main.py

[program:panel]
command=python panel/app.py
